# infra/systemd/robothor-app-pf.service
[Unit]
Description=Genus OS Helm Dashboard — Secondary Instance
After=network-online.target robothor-engine-pf.service redis.service
Wants=network-online.target

[Service]
User=robothor
Group=robothor
EnvironmentFile=/etc/robothor/robothor.env
ExecStartPre=${ROBOTHOR_WORKSPACE}/scripts/decrypt-secrets.sh
EnvironmentFile=-/run/robothor/secrets.env
ExecStart=/usr/bin/node ${ROBOTHOR_WORKSPACE}/app/.next/standalone/server.js
WorkingDirectory=${ROBOTHOR_WORKSPACE}/app
Restart=always
RestartSec=5
KillMode=control-group
Environment=NODE_ENV=production
Environment=PORT=3004
Environment=ROBOTHOR_ENGINE_URL=http://127.0.0.1:18800
Environment=AGENT_SESSION_KEY=agent:pf-helm:primary
Environment=HELM_AGENT_ID=pf-helm
Environment=REDIS_URL=redis://localhost:6379/0

[Install]
WantedBy=multi-user.target
