# scripts/bridge-watchdog.service
[Unit]
Description=Bridge Watchdog — health check and self-healing
After=robothor-bridge.service

[Service]
Type=oneshot
EnvironmentFile=/etc/robothor/robothor.env
ExecStart=${ROBOTHOR_WORKSPACE}/brain/memory_system/venv/bin/python ${ROBOTHOR_WORKSPACE}/brain/scripts/bridge_watchdog.py
User=robothor
WorkingDirectory=${ROBOTHOR_WORKSPACE}/brain
